# enhancements/gromox-requeue/gromox-requeue.path
[Unit]
Description=Gromox Requeue

[Path]
DirectoryNotEmpty=/var/lib/gromox/queue/save

[Install]
WantedBy=timers.target
