# enhancements/gromox-requeue/gromox-requeue.timer
[Unit]
Description=Gromox Requeue

[Timer]
OnCalendar=*:5/30
AccuracySec=1m
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
